# rootfs/etc/services.d/filebot-info/run
#!/bin/sh
exec /opt/filebot/filebot -script fn:sysinfo
